# bot.conf
[program:ivent_bot]
directory=/home/ivent_bot
command=python3 app.py
autostart=true
autorestart=true
stderr_logfile=/home/ivent_bot/logfile_err.log
stdout_logfile=/home/ivent_bot/logfile.log
